# 99-amdgpu-fan.rules
SUBSYSTEM=="hwmon", KERNEL=="hwmon*", ATTR{name}=="amdgpu", ACTION=="add", RUN+="/bin/chmod 666 /sys%p/pwm1 /sys%p/pwm1_enable"
